# systemd/camera.service
[Unit]
Description=Pi Camera

[Service]
Type=simple
ExecStart=/home/pi/rat/camera.py

[Install]
WantedBy=multi-user.target