# systemd/light_control.service
[Unit]
Description=Light Control

[Service]
Type=simple
ExecStart=/home/pi/rat/light_control.py

[Install]
WantedBy=multi-user.target