# systemd/uploader.service
[Unit]
Description=YouTube Upload Service

[Service]
Type=simple
ExecStart=/home/pi/rat/uploader.py

[Install]
WantedBy=multi-user.target